# packaging/deb/sendspin-service.service
[Unit]
Description=Sendspin playback service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/sendspin-service
ExecStart=/opt/sendspin-service/bin/sendspin-service
Restart=on-failure
RestartSec=2
DynamicUser=yes
StateDirectory=sendspin-service

[Install]
WantedBy=multi-user.target
